# data_enchilada/slave/templates/etc/supervisor/conf.d/data_enchilada.conf.erb
[program:data_enchilada]
command=java -jar /usr/lib/data_enchilada/data_enchilada.jar
autostart=true
autorestart=true
stderr_logfile=/var/log/data_enchilada.err.log
stdout_logfile=/var/log/data_enchilada.out.log