# enchilada/slave/templates/etc/supervisor/conf.d/enchilada.conf.erb
[program:enchilada]
command=java -jar /usr/lib/enchilada/enchilada.jar
autostart=true
autorestart=true
stderr_logfile=/var/log/enchilada.err.log
stdout_logfile=/var/log/enchilada.out.log