# enchilada/slave/templates/etc/supervisor/conf.d/enchilada.conf.erb
[program:enchilada]
command=java -jar /usr/lib/enchilada/enchilada.jar
autostart=true
autorestart=true
stderr_logfile=/var/log/enchilada.err.log
stdout_logfile=/var/log/enchilada.out.log

[program:ruby_enchilada]
command=bundle exec puma -p 3000 -e production
autostart=true
autorestart=true
stderr_logfile=/var/log/ruby_enchilada.err.log
stdout_logfile=/var/log/ruby_enchilada.out.log

